# ansible_web-api.service
[Unit]
Description=web-api

[Service]
Type=simple
User=root
Group=root
ExecStart=/etc/api/main
Restart=always

[Install]
WantedBy=multi-user.target
